# systemd/reconf_dragon_online.service
[Unit]
Description=Reconfigure Cluster with the new ID

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/reconf_dragon_online.py
